# scripts/delta_imports.crontab
*/1     *       *       *       *       sh /export/home/apps/search/search/scripts/product_deltaimport.sh >> /dev/null 2>&1